# committer.cron
# Sincroniza de forma sequencial cada repositório

*/30 * 	* * *	root committer
